# data/SecureFileAccess.service.in
[Unit]
Description=Secure File Access Daemon
After=syslog.target
Requires=dbus.service

[Service]
Type=dbus
BusName=org.mupuf.SecureFileAccess
ExecStart=@bindir@/sfad
KillMode=process
Restart=always
WatchdogSec=10
Environment=DISPLAY=:0


[Install]
WantedBy=graphical.target
Alias=dbus-org.mupuf.SecureFileAccess.service

